# infra/systemd/quickclip-api.service
[Unit]
Description=QuickClip Upload API Service
After=network.target

[Service]
Type=simple
User=quickclip
WorkingDirectory=/opt/quickclip/upload-api
EnvironmentFile=/etc/quickclip/env
ExecStart=/usr/bin/node server-vm.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
